# deploy/launchd/docker-watchdog.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label (must be unique) -->
    <key>Label</key>
    <string>com.llmclass.platform.watchdog</string>
    
    <!-- Program and Arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/Documents/llmclass_platform/deploy/launchd/watchdog.sh</string>
    </array>
    
    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/Documents/llmclass_platform/deploy</string>
    
    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>DOCKER_HOST</key>
        <string>unix:///Users/YOUR_USERNAME/.docker/run/docker.sock</string>
    </dict>
    
    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>
    
    <!-- Auto-start settings -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Don't restart on successful exit -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Documents/llmclass_platform/logs/watchdog.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Documents/llmclass_platform/logs/watchdog.error.log</string>
    
    <!-- Process settings -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- Security (run as user, not root) -->
    <key>UserName</key>
    <string>YOUR_USERNAME</string>
    
    <key>GroupName</key>
    <string>staff</string>
    
    <!-- Additional settings -->
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>

<!-- 
This watchdog service runs every 5 minutes to:
1. Check if Docker services are running
2. Restart failed containers
3. Monitor health endpoints
4. Log status information

Installation:
1. First install watchdog.sh script
2. Copy this file to ~/Library/LaunchAgents/com.llmclass.platform.watchdog.plist
3. Edit paths and username
4. Load with: launchctl load ~/Library/LaunchAgents/com.llmclass.platform.watchdog.plist
-->