# config/launchd/com.llmclassroom.unified.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.llmclassroom.unified</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jjh_server/llmclass_platform/llmclass</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/jjh_server/llmclass_platform</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
        <key>NetworkState</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/jjh_server/llmclass_platform/logs/unified.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jjh_server/llmclass_platform/logs/unified.error.log</string>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>